# flask-mega-downloader.service
[Unit]
Description=Flask Mega Downloader
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/flask-mega-downloader
Environment=MEGA_DOWNLOADER_CONFIG=/etc/flask-mega-downloader/config.py
Environment=HOME=/var/www
ExecStart=/opt/flask-mega-downloader/.venv/bin/python /opt/flask-mega-downloader/server.py
Restart=on-failure
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
